# deploy/systemd/investment-dashboard-tunnel.service
[Unit]
Description=Investment Dashboard Cloudflare Tunnel
Wants=network-online.target investment-dashboard-web.service
After=network-online.target investment-dashboard-web.service
StartLimitIntervalSec=120
StartLimitBurst=10

[Service]
Type=simple
User=ericarthuang
WorkingDirectory=/home/ericarthuang/.openclaw/workspace/investment_dashboard
ExecStart=/home/ericarthuang/.openclaw/workspace/investment_dashboard/deploy/scripts/run_tunnel.sh
Restart=always
RestartSec=5
TimeoutStopSec=20
KillSignal=SIGINT
Environment=PYTHONUNBUFFERED=1
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
